# kibana.conf
[program:kibana]
command=/opt/kibana/bin/kibana
user=kibana
redirect_stderr=true
directory=/opt/kibana