# kibana.conf
[program:kibana]
command=/opt/kibana/bin/kibana --max-old-space-size=200
user=kibana
redirect_stderr=true
directory=/opt/kibana